# share/xmz-mod-touch-server.service
[Unit]
Description="LED und RELASI Control"
After=multi-user.target

[Service]
ExecStart=/root/server &

[Install]
WantedBy=multi-user.target
